# states/os/quantum/etc/quantum-dhcp-agent.upstart
start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start \
                       --chuid quantum \
                       --exec /usr/local/bin/quantum-dhcp-agent \
                       -- --config-dir /etc/quantum \
                          --config-file /etc/quantum/quantum.conf \
                          --config-file /etc/quantum/dhcp_agent.ini
