# states/os/keystone/keystone.upstart
start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start \
                       --chuid keystone \
                       --exec /usr/local/bin/keystone-all \
                       -- --config-dir /etc/keystone
